# yeelight-smartthings.service
[Unit]
Description=node.js server to manage yeelight bulbs

[Service]
ExecStart=/var/opt/node/yeelight-smartthings/server.js
Restart=always
User=nobody
Group=nobody
Environment=PATH=/usr/bin:/usr/local/bin
Environment=NODE_ENV=production
WorkingDirectory=/var/opt/node/yeelight-smartthings/

[Install]
WantedBy=multi-user.target